# laravel-job-runner/docker/worker/supervisord.conf
[supervisord]
nodaemon=true

[program:worker1]
command=php worker/worker.php
autostart=true
autorestart=true

[program:worker2]
command=php worker/worker.php
autostart=true
autorestart=true
